# rockchip/etc/crontabs/root
*/15 * * * * /etc/init.d/ddns restart #定时重启DDNS
56 23 * * 6 docker restart nextcloud #定时重启nextcloud 

