# rockchip/etc/crontabs/root
56 23 * * * docker restart nextcloud #定时重启nextcloud 

