# rockchip/etc/crontabs/root
56 23 * * 6 docker restart nextcloud #定时重启nextcloud 

